# supervisord/supervisor_services.conf
[program:nginx]
commadn=/usr/sbin/nginx -g 'deamon off;'
autostart=true
autorestart=true
